# audino/mysql/crontab
*/1 * * * * cp /var/lib/mysql/* /mnt/backup >> /var/log/cron.log 2>&1
