# audino/mysql/crontab
0 */6 * * * mysqldump -u root --password=root audino > /mnt/backup/`date +\%Y\%m\%d-\%H:\%M`.sql 
#emptyline required for cron to run properly
